# weeecycle.service
[Unit]
Description=Weeecycle Node.js Server
After=network.target

[Service]
Type=simple
User=swhitelex
WorkingDirectory=/home/swhitelex/weeecycle
ExecStart=/usr/bin/node /home/swhitelex/weeecycle/server.js
Restart=on-failure
# Environment variables can be added here
# Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
